# roles/common/files/logwatch.service
[Unit]
Description=Logwatch generate raport

[Service]
Type=oneshot
ExecStart=sh -c 'logwatch | tee /tmp/logwatch.log'
